# turnkey-owntracks/overlay/etc/supervisord.d/m2s.ini
[program:m2s]
command = /usr/local/owntracks/m2s/m2s-launcher.sh
directory = /usr/local/owntracks/m2s
user = mosquitto
# environment= LOGNAME=jpm
