# turnkey-owntracks/overlay/etc/supervisord.d/wss.ini
[program:WSS]
command = /usr/local/owntracks/wss/WSS --brokerHost 127.0.0.1 --brokerPort 1883 --websocketPort 18883
directory = /usr/local/owntracks/wss
user = mosquitto
